# roles/sshpiper/templates/sshpiperd.service.j2
[Unit]
Description=SSHPiper Service
After=network.target

[Service]
Type=simple
ExecStart= {{ sshpiper_bin_dir }}/sshpiperd -p 22 --log-level=trace {{ sshpiper_bin_dir }}/yaml --config {{ sshpiper_dest_dir }}/sshpiperd.yaml --no-check-perm
Restart=on-failure
User=root

[Install]
WantedBy=multi-user.target
